# roles/cs.mageops-cli/files/mageops-update-features.service
[Unit]
Description=Update dynamic features

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mageopscli apply_features
