# roles/cs.mageops-cli/files/mageops-update-features.timer
[Unit]
Description=Trigger update of dynamic features

[Timer]
OnActiveSec=10sec
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
