# roles/cs.php-opcache-prewarm/templates/preload-magento-opcache.timer
[Unit]
Description=Run opcache preload on boot

[Timer]
OnBootSec=20s

[Install]
WantedBy=timers.target
