# roles/cs.mageops-cli-profile/templates/profile.d/001-mageops-locale.sh
export LANG={{ mageops_lang }}.UTF-8
export LANGUAGE={{ mageops_lang }}.UTF-8
export LC_COLLATE=C
export LC_CTYPE={{ mageops_lang }}.UTF-8